# systemd/broken_service.service
[Unit]
Description=Broken HPC Service
After=network.target

[Service]
ExecStart=/opt/app/run.sh
Restart=always

[Install]
WantedBy=multi-user.target
